# systemd/logger.service
[Unit]
Description=Vehicle OBD + GPS Logger
After=bluetooth.target gpsd.service rfcomm-bind.service
Wants=bluetooth.target gpsd.service rfcomm-bind.service

[Service]
ExecStart=/usr/bin/python3 /home/nero/vehicle-logger/logger.py
WorkingDirectory=/home/nero/vehicle-logger
Restart=on-failure
RestartSec=5
User=nero
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
